# clearpath_robot/services/clearpath-vcan.service
[Unit]
Description="Clearpath Virtual CAN over UDP Bridge"
After=network.target

[Service]
Restart=on-failure
User=root
Type=forking
ExecStart=/bin/bash -e /etc/clearpath/vcan-start

[Install]
WantedBy=multi-user.target
